# scripts/podman/propai.container.in
# PropAi Sync gateway — Podman Quadlet (rootless)
# Installed by setup-podman.sh into PropAi Sync's ~/.config/containers/systemd/
# {{PROPAI_HOME}} is replaced at install time.

[Unit]
Description=PropAi Sync gateway (rootless Podman)

[Container]
Image=propai:local
ContainerName=propai
UserNS=keep-id
# Keep container UID/GID aligned with the invoking user so mounted config is readable.
User=%U:%G
Volume={{PROPAI_HOME}}/.propai:/home/node/.propai:Z
EnvironmentFile={{PROPAI_HOME}}/.propai/.env
Environment=HOME=/home/node
Environment=TERM=xterm-256color
PublishPort=18789:18789
PublishPort=18790:18790
Pull=never
Exec=node dist/index.js gateway --bind lan --port 18789

[Service]
TimeoutStartSec=300
Restart=on-failure

[Install]
WantedBy=default.target



